# dot_local/share/systemd/user/ledger-prices.timer
[Unit]
Description=Run ledger-prices on boot and every one minute

[Timer]
OnBootSec=15min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
