# dot_local/share/systemd/user/ledger-prices-commit.timer
[Unit]
Description=Run git commit/push every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=default.target
