# dot_local/share/systemd/user/torn.service
[Unit]
Description=Update torn.json

[Service]
ExecStart=%h/scripts/torn --update

[Install]
WantedBy=default.target
