# dot_local/share/systemd/user/torn.timer
[Unit]
Description=update torn.json

[Timer]
OnBootSec=1sec
OnUnitActiveSec=10sec

[Install]
WantedBy=timers.target
